# deploy/dora-proxy.service
[Unit]
Description=DORA Jira CORS Proxy
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/dora
EnvironmentFile=/home/ubuntu/dora/.env
ExecStart=/home/ubuntu/dora/venv/bin/python3 proxy/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
